# crontab
* * * * * cd ~/BelPexPLCnext && /usr/bin/python3 app.py --user admin --password ??????